# ops/etc/supervisor/conf.d/laravel-worker-ainicheng-echo-server.conf
[program:laravel-worker-ainicheng-echo-server]
process_name=%(program_name)s_%(process_num)02d
directory=/data/www/ainicheng.com
command=/data/build/node-v10.16.0-linux-x64/bin/laravel-echo-server start
autostart=true
autorestart=true
user=www
numprocs=1
redirect_stderr=true
stdout_logfile=/data/www/ainicheng.com/storage/logs/worker.log