# ops/workers/web/laravel-worker-ainicheng-echo-server.conf
[program:laravel-worker-ainicheng-echo-server]
process_name=%(program_name)s_%(process_num)02d
directory=/data/www/ainicheng.com
command=/data/.nvm/versions/node/v10.16.0/bin/node /data/.nvm/versions/node/v10.16.0/bin/laravel-echo-server start
autostart=true
autorestart=true
user=www
numprocs=1
redirect_stderr=true
stdout_logfile=/data/www/ainicheng.com/storage/logs/echo-server.log
